# files/etc/uci-defaults/99-setup-modem
#!/bin/sh
uci set semutcriminals.modem.apn='internet'
uci commit semutcriminals
