# files/etc/rc.local
#!/bin/sh -e
# rc.local

# Aktifkan modem + ping test
/etc/init.d/modemctl start
ping -c 4 8.8.8.8 > /tmp/pingtest.txt &

exit 0
